# root/etc/services.d/apache/run
#!/usr/bin/execlineb -P
httpd -DFOREGROUND;
