# root/etc/services.d/mysql/run
#!/usr/bin/execlineb -P
s6-setuidgid mysql
mysqld --init-file=/scripts/mysql-init --user=mysql --console
